# formtopia-cookbook/templates/default/formtopia-inspector.conf.erb
description "Formtopia Inspector"
author "Chef"

start on (filesystem and net-device-up)
stop on runlevel [!2345]

respawn
respawn limit 5 30

chdir <%= @dir %>
setuid <%= @user %>



script
  NODE_ENV=production node-inspector
end script

emits formtopia-inspector-running